# .crontab
# ~/.crontab: scheduled tasks via crontab. Executed by cron.
# Use crontab -e to add tasks below. Use crontab -l to list them.

# Copy containers in ~/images to external hard drive (mounted as d by wsl).
0 16 * * * rsync --recursive --verbose $HOME/images /mnt/d/backups --log-file=$HOME/logs/rsync-$(date +"%Y-%m-%d").log
